# install/operating_systems_laboratories_4.service
[Unit]
Description=Run monitor

[Service]
Type=oneshot
ExecStart=/usr/share/operating_systems_laboratories_4/bin/monitor